# system_files/homeserver/etc/systemd/system/bootc-system-setup.service
[Unit]
Description=Configure system
Wants=network-online.target
After=network-online.target systemd-sysusers.service systemd-tmpfiles-setup.service 
Requires=systemd-sysusers.service systemd-tmpfiles-setup.service 

[Service]
Type=oneshot
ExecStart=/usr/libexec/system-setup
User=root
RemainAfterExit=yes
Environment=TARGET_USER=caspertdk
Environment=TARGET_ID=1010

[Install]
WantedBy=multi-user.target
