# system_files/homeserver/etc/systemd/system/bootc-system-setup.service
[Unit]
Description=Configure system
#After=rpm-ostreed.service
#Before=systemd-user-sessions.service
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/libexec/system-setup
User=root
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target