# system_files/homeserver-gui/usr/lib/systemd/system/bootc-fetch-apply-updates.service.d/override.conf
# we dont want to reboot the system after applying updates, as its probably a client machine (laptop) and not a server
[Service]
ExecStart=
ExecStart=/usr/bin/bootc upgrade --quiet