# system_files/homeserver-gui/usr/lib/systemd/system/bootc-fetch-apply-updates.timer.d/updates.conf
[Timer]
# Clear the original triggers
OnBootSec=
OnUnitInactiveSec=
RandomizedDelaySec=

# New calendar-based schedule: every day at 12:30
OnCalendar=*-*-* 12:30:00
Persistent=true